# etc/systemd/pwmfan-controller.service
[Unit]
Description=PWM Fan Controller Service
Requires=pwmfan-setup.service
After=multi-user.target pwmfan-setup.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/pwmfan_controller.py --mode auto --verbose
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
